# rootfs/etc/cron.d/cron.hourly/wp-autoupdate
#!/usr/bin/env bash

echo $(date) >> /tmp/autoupdate
bash /xshok-wp-autoupdate.sh >> /tmp/autoupdate
